# systemd/kayak-sync.timer
[Unit]
Description=Run kayak observation sync twice daily

[Timer]
OnCalendar=06:15
OnCalendar=18:15
Persistent=true

[Install]
WantedBy=timers.target
