# systemd/kayak-pipeline.service
[Unit]
Description=Kayak river levels pipeline
After=network-online.target

[Service]
Type=oneshot
User=pat
WorkingDirectory=/home/pat/kayak
EnvironmentFile=/home/pat/.config/kayak/.env
ExecStart=/home/pat/.venv/bin/levels pipeline
TimeoutStartSec=300
Nice=10
