# systemd/kayak-notify-failure@.service
[Unit]
Description=Send failure notification for %i

[Service]
Type=oneshot
User=pat
# `-` prefix: unit still works if the env file is missing; the ntfy branch
# in ExecStart is guarded on NTFY_TOPIC being non-empty, so a missing file
# just disables the push channel (email + syslog paths unaffected).
# Paths in WorkingDirectory= / EnvironmentFile= / ReadWritePaths= AND
# the binary path in ExecStart= itself are LITERAL — systemd does not
# expand $VAR in those positions (executable name must be an absolute
# literal per systemd.exec). ExecStart= arguments and ExecStartPost=
# DO expand ${KAYAK_HOME} (set by /etc/kayak/env with the in-unit
# fallback default below).
Environment=KAYAK_HOME=/home/pat
EnvironmentFile=-/etc/kayak/env
EnvironmentFile=-/home/pat/.config/kayak/.env
# Log to syslog (source of truth) then attempt email + push. Each channel
# pipes any failure back to syslog so broken msmtp or curl doesn't silently
# drop alerts. Requires /etc/msmtprc to have a real Gmail app password set;
# NTFY_TOPIC must be set in the EnvironmentFile for the push channel.
ExecStart=/bin/bash -c 'MSG="ALERT: %i failed at $(date -u +\"%%Y-%%m-%%dT%%H:%%M:%%SZ\")"; echo "$MSG" | logger -t kayak-alert -p user.err; echo "$MSG" | mail -s "Kayak: %i failed" pat.kayak@gmail.com 2>&1 | logger -t kayak-alert -p user.err; if [ -n "${NTFY_TOPIC:-}" ]; then echo "$MSG" | curl -fsS -m 10 --retry 3 -d @- "ntfy.sh/$NTFY_TOPIC" -H "Title: Kayak: %i failed" -H "Priority: high" 2>&1 | logger -t kayak-alert -p user.err; fi'
