# systemd/kayak-cert-renewal-test.service
[Unit]
Description=Kayak: weekly certbot state sanity check
After=network-online.target
Wants=network-online.target
OnFailure=kayak-notify-failure@%n.service

[Service]
Type=oneshot
User=root
# Paths in WorkingDirectory= / EnvironmentFile= / ReadWritePaths= AND
# the binary path in ExecStart= itself are LITERAL — systemd does not
# expand $VAR in those positions (executable name must be an absolute
# literal per systemd.exec). ExecStart= arguments and ExecStartPost=
# DO expand ${KAYAK_HOME} (set by /etc/kayak/env with the in-unit
# fallback default below).
Environment=KAYAK_HOME=/home/pat
EnvironmentFile=-/etc/kayak/env
EnvironmentFile=/home/pat/.config/kayak/.env
# `certbot certificates` is read-only: parses /etc/letsencrypt/renewal/*.conf,
# inspects each lineage's fullchain.pem, prints status. No ACME traffic, no
# nginx mutation, no temp_checkpoint. Replaced `renew --dry-run` on 2026-05-18
# after the --nginx plugin's revert step failed and left acme-challenge stubs
# in /etc/nginx/sites-available/. Real renewals still run via the system
# certbot.timer twice daily; kayak-cert-expiry.timer catches upcoming expiry
# independently. Unit filename predates the purpose change.
ExecStart=/usr/bin/certbot certificates
ExecStartPost=-/usr/bin/curl -fsS -m 10 --retry 3 -o /dev/null ${HC_CERT_RENEWAL_TEST}
TimeoutStartSec=60
Nice=15

# certbot certificates needs:
#   - read+write /etc/letsencrypt — even the read-only `certificates`
#     subcommand acquires .certbot.lock in the config dir on startup
#   - write /var/log/letsencrypt (logs to letsencrypt.log on every invocation)
#   - write /var/lib/letsencrypt (OCSP cache, etc.)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/etc/letsencrypt /var/log/letsencrypt /var/lib/letsencrypt
ProtectHome=read-only
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
# review-4 R5.3: match the User=pat units' syscall hardening. certbot keeps its
# default root caps (CapabilityBoundingSet left unset -- it needs file caps for
# /etc/letsencrypt); @system-service is the standard service syscall allowlist.
SystemCallFilter=@system-service
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM
UMask=0022
