# systemd/kayak-fetch-osmb.service
[Unit]
Description=Fetch Oregon SMB hazard/access GeoJSON overlays
After=network-online.target
OnFailure=kayak-notify-failure@%n.service

[Service]
Type=oneshot
User=pat
WorkingDirectory=/home/pat/kayak
# Paths in WorkingDirectory= / EnvironmentFile= / ReadWritePaths= AND
# the binary path in ExecStart= itself are LITERAL — systemd does not
# expand $VAR in those positions (executable name must be an absolute
# literal per systemd.exec). ExecStart= arguments and ExecStartPost=
# DO expand ${KAYAK_HOME} (set by /etc/kayak/env with the in-unit
# fallback default below).
Environment=KAYAK_HOME=/home/pat
EnvironmentFile=-/etc/kayak/env
EnvironmentFile=/home/pat/.config/kayak/.env
ExecStart=/home/pat/.venv/bin/levels fetch-osmb
# Healthchecks.io heartbeat. Leading `-` so a ping failure (HC down,
# DNS hiccup) doesn't mark the whole service failed — the fetch itself
# already succeeded by this point. HC_FETCH_OSMB lives alongside the
# other HC_* URLs in /home/pat/.config/kayak/.env (loaded above).
ExecStartPost=-/usr/bin/curl -fsS -m 10 --retry 3 -o /dev/null ${HC_FETCH_OSMB}
TimeoutStartSec=300
Nice=15

# Sandboxing. The fetcher only writes GeoJSON files to the OSMB staging dir
# (OSMB_DIR; default BASE_DIR/var/osmb — `levels build` later copies them into
# OUTPUT_DIR/static) and reads no other state — no DB, no secrets.
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=read-only
# OSMB_DIR default is /home/pat/kayak/var/osmb. The dir is kept present by the
# tracked var/osmb/.gitkeep so systemd's setup-time existence check passes on a
# fresh clone; scope write to just it (not the whole var/, which holds other
# components' runtime state).
ReadWritePaths=/home/pat/kayak/var/osmb
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
CapabilityBoundingSet=
SystemCallFilter=@system-service
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM
UMask=0077
