# systemd/kayak-decimate.timer
[Unit]
Description=Run kayak decimate daily

[Timer]
OnCalendar=02:32
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
