# systemd/kayak-status.timer
[Unit]
Description=Nightly render of the operator status page

[Timer]
# 03:30 — one hour before kayak-metadata-snapshot.timer at 04:30, so the
# status page reflects the same on-disk state the snapshot then captures.
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
