# systemd/kayak-audit-gauges.timer
[Unit]
Description=Audit gauge metadata twice a month (2nd + 17th)

[Timer]
# OnCalendar gives stable 03:00 fire times. Gaps alternate ~15 days
# (2->17) and 13-16 days (17->2 of next month, varies by month length).
# The matching service uses --days 16 so the audit's lookback window
# always meets or exceeds the max calendar gap — every stop/start event
# is caught by at least one fire, with overlap in shorter-gap months.
OnCalendar=*-*-02,17 03:00
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
