# systemd/kayak-cert-renewal-test.service
[Unit]
Description=Kayak: weekly certbot renew --dry-run
After=network-online.target
Wants=network-online.target
OnFailure=kayak-notify-failure@%n.service

[Service]
Type=oneshot
User=root
EnvironmentFile=/home/pat/.config/kayak/.env
ExecStart=/usr/bin/certbot renew --dry-run --quiet --non-interactive --no-random-sleep-on-renew
ExecStartPost=-/usr/bin/curl -fsS -m 10 --retry 3 -o /dev/null ${HC_CERT_RENEWAL_TEST}
TimeoutStartSec=300
Nice=15

# certbot --nginx --dry-run needs to:
#   - read /etc/letsencrypt (cert state)
#   - write /var/log/letsencrypt + /var/lib/letsencrypt (certbot work)
#   - run `nginx -t`, which opens /var/log/nginx/error.log for append
#   - have the --nginx plugin insert/remove a temporary location block in
#     /etc/nginx/ during ACME validation, even on dry-run
#   - reload nginx via the pid in /run/nginx.pid
# Plus /home/pat/.config/kayak/.env for HC_CERT_RENEWAL_TEST. ProtectHome
# is read-only (matches every other kayak-* unit's pattern) so root can still
# read the env file while /home is otherwise locked.
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/etc/letsencrypt /var/log/letsencrypt /var/lib/letsencrypt /etc/nginx /var/log/nginx /run
ProtectHome=read-only
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
UMask=0022
