# systemd/kayak-sync.service
[Unit]
Description=Kayak observation sync from legacy MySQL
After=network-online.target

[Service]
Type=oneshot
User=pat
WorkingDirectory=%h/kayak
EnvironmentFile=%h/.config/kayak/.env
ExecStart=%h/kayak/systemd/kayak-sync.sh
TimeoutStartSec=600
Nice=10
