# systemd/kayak-heartbeat.timer
[Unit]
Description=Run kayak heartbeat weekly

[Timer]
OnCalendar=Sun 06:00
Persistent=true

[Install]
WantedBy=timers.target
