# systemd/kayak-status.service
[Unit]
Description=Render the operator status page (nightly)
OnFailure=kayak-notify-failure@%n.service

[Service]
Type=oneshot
User=pat
WorkingDirectory=/home/pat/kayak
# Paths in WorkingDirectory= / EnvironmentFile= / ReadWritePaths= AND
# the binary path in ExecStart= itself are LITERAL — systemd does not
# expand $VAR in those positions (executable name must be an absolute
# literal per systemd.exec). ExecStart= arguments and ExecStartPost=
# DO expand ${KAYAK_HOME} (set by /etc/kayak/env with the in-unit
# fallback default below).
Environment=KAYAK_HOME=/home/pat
EnvironmentFile=-/etc/kayak/env
EnvironmentFile=-/home/pat/.config/kayak/.env
ExecStart=/home/pat/.venv/bin/levels status --output /home/pat/kayak/var/status.html
# `-` prefix tolerates an empty/missing HC_STATUS env var — same convention
# as kayak-healthcheck.service. Create the check at healthchecks.io and
# paste the ping URL into ~/.config/kayak/.env to opt in.
ExecStartPost=-/usr/bin/curl -fsS -m 10 --retry 3 -o /dev/null ${HC_STATUS}
TimeoutStartSec=120
Nice=15

# Sandboxing — modeled on kayak-recap.service (read journald + write a
# single file). Differences:
#   * ReadWritePaths=/home/pat/kayak/var so the rendered HTML can land.
#   * ProcSubset=all (not =pid): the disk-and-memory section reads
#     /proc/meminfo, which the pid-only subset hides.
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/pat/kayak/var
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=all
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
CapabilityBoundingSet=
SystemCallFilter=@system-service
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM
UMask=0077
