# systemd/kayak-metadata-snapshot.timer
[Unit]
Description=Daily kayak metadata snapshot to git

[Timer]
OnCalendar=*-*-* 04:30:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
