# systemd/kayak-status.timer
[Unit]
Description=Nightly render of the operator status page

[Timer]
# 03:30 — a quiet pre-dawn slot for the nightly status render, clear of the
# pipeline/backup timers.
OnCalendar=*-*-* 03:30:00
RandomizedDelaySec=300
Persistent=true

[Install]
WantedBy=timers.target
